# flight-tracker.service
[Unit]
Description=E-Ink Flight Tracker
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/flight_tracker
ExecStart=/home/pi/flight_tracker/venv/bin/python3 main.py
Restart=on-failure
RestartSec=30
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
